# deploy/spannora.service
[Unit]
Description=spannora — self-hosted web chat for Claude Code
Documentation=https://github.com/gididaf/spannora
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/spannora

# Run the compiled JS with the system node (>= 20).
ExecStart=/usr/bin/env node /opt/spannora/dist/server.js

Restart=on-failure
RestartSec=5

# State dir → /var/lib/spannora. Owned 0700 by the service user (root by
# default). SQLite DB lives here so it survives package upgrades.
StateDirectory=spannora
StateDirectoryMode=0700

Environment=NODE_ENV=production
Environment=SPANNORA_HOST=127.0.0.1
Environment=SPANNORA_PORT=7878
Environment=SPANNORA_DB=/var/lib/spannora/spannora.db
Environment=IS_SANDBOX=1
# Optional retention. Uncomment + tune to auto-prune old conversations
# (deletes the matching ~/.claude/projects/.../<session-id>.jsonl too).
# Environment=SPANNORA_RETENTION_DAYS=30

# No User= → systemd defaults to root. This matches the typical
# Claude-on-a-VM workflow: the SDK reuses /root/.claude/ for auth, and
# Claude's tool calls inherit root's filesystem access. If you'd rather
# sandbox spannora to a less-privileged user, add `User=` and `Group=`
# lines and re-run the installer.
PrivateTmp=true

[Install]
WantedBy=multi-user.target
